# init/krypt-daemon.service
[Unit]
Description=Krypt VM Management Daemon
Documentation=https://github.com/sudoshaco/krypt-os
# Xen-Store muss erreichbar sein bevor der Daemon VM-Operationen ausführt.
# Falls xenstore.service nicht existiert (Test-Umgebung) greift After= nicht.
After=local-fs.target xenstore.service
Wants=xenstore.service

[Service]
Type=simple
User=root
Group=root

ExecStart=/usr/bin/krypt-daemon
ExecReload=/bin/kill -HUP $MAINPID

Restart=on-failure
RestartSec=5s
# Nie neu starten wenn absichtlich gestoppt (SIGTERM/SIGINT)
RestartPreventExitStatus=0

# /run/krypt/ anlegen (IPC-Socket, boot-stick-serial)
RuntimeDirectory=krypt
RuntimeDirectoryMode=0700

# /etc/krypt/ für daemon.toml und panic.toml
ConfigurationDirectory=krypt
ConfigurationDirectoryMode=0700

# Sicherheitshärtung — soweit mit xl-Calls kompatibel
# NoNewPrivileges=yes ist nicht möglich: xl nutzt setuid-Bits
ProtectSystem=strict
ProtectHome=yes
# Daemon liest /etc/krypt/daemon.toml read-only — keine Write-Notwendigkeit.
# RuntimeDirectory= erzeugt /run/krypt; nur dort wird der IPC-Socket gebunden
# und die boot-stick-serial Datei abgelegt. ReadWritePaths bewusst minimal.
ReadWritePaths=/run/krypt
PrivateTmp=yes
# xl braucht /dev/xen — Inline-Kommentare nach "=value" parst systemd NICHT
# als Kommentar, sondern als Teil des Werts → "Invalid argument".
PrivateDevices=no
ProtectKernelTunables=yes
ProtectControlGroups=yes
# xl-Binaries sind setuid root — RestrictSUIDSGID=yes würde xl create blockieren.
RestrictSUIDSGID=no
LockPersonality=yes
RestrictRealtime=yes

StandardOutput=journal
StandardError=journal
SyslogIdentifier=krypt-daemon

[Install]
WantedBy=multi-user.target
