# init/krypt-boot-agent.service
[Unit]
Description=Krypt Boot Auth Stick Registration
Documentation=https://github.com/sudoshaco/krypt-os
# Läuft einmalig nach krypt-daemon — registriert den beim initramfs-Boot
# verwendeten Auth-Stick so dass der USB-Monitor ihn gezielt überwacht.
After=krypt-daemon.service
Requires=krypt-daemon.service
# Nur relevant wenn ein LUKS-Root-Device tatsächlich via krypt-Hook geöffnet wurde
ConditionPathExists=/dev/mapper/krypt-root

[Service]
Type=oneshot
User=root
Group=root
# Script liegt unter /usr/lib/krypt/ (kein $PATH-Scanning nötig)
ExecStart=/usr/lib/krypt/krypt-boot-agent.sh
# Service gilt als "aktiv" nach einmaligem Durchlauf (für Abhängigkeitsketten)
RemainAfterExit=yes

StandardOutput=journal
StandardError=journal
SyslogIdentifier=krypt-boot-agent

# Sicherheitshärtung — analog zu krypt-daemon, soweit für ein oneshot-Skript
# das dmsetup/lsblk/udevadm aufruft sinnvoll ist.
# Boot-Agent braucht NIE Privilege-Escalation; setuid-Binaries sind hier
# kein Issue (im Unterschied zum daemon, der xl aufrufen muss).
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=yes
ReadWritePaths=/run/krypt
PrivateTmp=yes
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectControlGroups=yes
RestrictSUIDSGID=yes
RestrictRealtime=yes
LockPersonality=yes

[Install]
WantedBy=multi-user.target
