# build/krypt-profile/airootfs/etc/systemd/system/krypt-installer-tty1.service
[Unit]
Description=Krypt OS TUI-Installer
After=multi-user.target
# Nicht starten wenn das System bereits installiert wurde
ConditionPathExists=!/etc/krypt/.installed

[Service]
Type=simple
ExecStart=/usr/bin/krypt-install
StandardInput=tty
StandardOutput=tty
StandardError=journal
TTYPath=/dev/tty1
TTYReset=yes
TTYVHangup=yes
Restart=on-failure
RestartSec=3s

[Install]
WantedBy=multi-user.target
